# dotconfig/systemd/user/audio.service
[Unit]
Description=Set Audio Volume

[Service]
Type=simple
Restart=Always
ExecStart=/home/caue/dotfiles/scripts/audio.fish
