# dotconfig/systemd/user/audio.service
[Unit]
Description=Set audio for 60%

[Service]
Type=simple
Restart=Always
ExecStart=pactl set-sink-volume @DEFAULT_SINK@ 60%
