# dotconfig/systemd/user/brightness.service
[Unit]
Description=Set screen brightness

[Service]
Type=oneshot
ExecStart=%h/.local/bin/set-brightness.sh
