# dotconfig/systemd/user/dotfiles.service
[Unit]
Description=Update dotfiles website

[Service]
Type=simple
ExecStart=/home/caue/Dev/go/api/s.fish
Restart=Always
