# dotconfig/systemd/user/limpante.service
[Unit]
Description=Clear RAM cache

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo /home/caue/dotfiles/scripts/limpante
