# dotconfig/systemd/brightness.service
[Unit]
Description=Set screen brightness

[Service]
Type=oneshot
ExecStart=/home/caue/.local/bin/set-brightness.sh
