# dotconfig/systemd/weekly-backup.service
[Unit]
Description=Weekly Backup

[Service]
Type=oneshot
ExecStart=/home/caue/dotfiles/scripts/weekly-backup.sh
