# dotconfig/systemd/weekly-backup.timer
[Unit]
Description=Run backup every Sunday

[Timer]
OnCalendar=Sun *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
